# distrod/libs/tests/resources/systemdunit/units/referenced_by_also2.service
[Unit]
Description=unit
Wants=network.target

[Service]
ExecStart=/bin/exit
Restart=always

[Install]
WantedBy=multi-user.target
