# distrod_packer/resources/run/systemd/system/systemd-tmpfiles-clean.service.d/exclude_wslg_sockets.conf
# This file is part of Distrod.
#
# Do not edit this file directly. This file can be overwritten when Distrod is updated.
# Instead, create /etc/systemd/system/systemd-tmpfiles-clean.service.d/exclude_wslg_sockets.conf, 
# systemd will load it.

# avoid `/tmp/.X11-unix` to be cleaned up by systemd
[Service]
ExecStart=
ExecStart=systemd-tmpfiles --clean --exclude-prefix=/tmp/.X11-unix